# revo/revo-roadrunner-mx7/systemd/kernel-cmdline.path
[Unit]
Description=Monitor /boot/cmdline.txt for changes
Documentation=https://github.com/revolution-robotics/roadrunner-debian

[Path]
PathChanged=/boot/cmdline.txt

[Install]
WantedBy=multi-user.target
